# app/supervisor/conf/001.conf
[program:torrentfs]
command=bash /root/app/anacrolix/torrentfs.sh
autostart=true
startsecs=10
autorestart=true
startretries=3
user=root
priority=999
redirect_stderr=true
stdout_logfile=/root/app/anacrolix/stdout.log
stopasgroup=false
killasgroup=false

[program:webdav]
command=/usr/bin/rclone serve webdav google: --read-only --addr :8001
autostart=true
startsecs=10
autorestart=true
startretries=3
user=root
priority=999
redirect_stderr=true
stdout_logfile=/root/app/rclone/webdav.log
stopasgroup=false
killasgroup=false
