# app/supervisor/conf/001.conf
[program:webdav]
command=/usr/bin/rclone serve webdav sharepoint2: --read-only --addr :8001
autostart=true
startsecs=10
autorestart=true
startretries=3
user=root
priority=999
redirect_stderr=true
stdout_logfile=/root/app/rclone/webdav.log
stopasgroup=false
killasgroup=false

[program:mount_google]
command=/root/app/rclone/mount.sh
autostart=true
startsecs=10
autorestart=true
startretries=3
user=root
priority=999
redirect_stderr=true
stdout_logfile=/root/app/rclone/mount_google.log
stopasgroup=false
killasgroup=false
